# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:memcached]
command=memcached -m 128 -u root
autorestart=true